# systemd-examples/feeds-to-bag.service
[Unit]
Description=Feeds to Bag

[Service]
Type=oneshot
ExecStart=feeds-to-bag ~/feeds-to-bag.yaml
